# data/export/ember_upstart/process.conf.erb
start on starting <%= app %>-<%= process.name %>
stop on stopping <%= app %>-<%= process.name %>
respawn

exec cd <%= engine.directory %> && <%= engine.environment.each_pair.map {|var, value| "export #{var.upcase}=#{value}; " }.join %><%= process.command %> >> <%= log_root %>/<%= process.name %>-<%= num %>.log 2>&1
