# data/export/ember_upstart/process.conf.erb
start on starting <%= app %>-<%= process.name %>
stop on stopping <%= app %>-<%= process.name %>
respawn

exec su - <%= user %> -c 'cd <%= engine.directory %>; <%= engine.environment.each_pair.map {|var,env| %{export #{var}="#{env}"; } }.join("") %> <%= process.command %> >> <%= log_root %>/<%= process.name %>.log 2>&1'
