# data/export/ember_upstart/process.conf.erb
start on starting <%= app %>-<%= process.name %>
stop on stopping <%= app %>-<%= process.name %>
respawn
kill signal QUIT
<%# Limit to 900MB soft / 1GB hard ram limit if it's resque process %>
<% if process.command[/resque/i] %>
limit rss 943718400 1073741824
<% end %>


script
  cd <%= engine.directory %>
<% engine.environment.each_pair do |var, value| %>
  export <%= var.upcase %>=<%= value %>
<% end %>
  <%= process.command %> >> <%= log_root %>/<%= process.name %>-<%= num %>.log 2>&1
end script
